# src/etc/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:transmission]
command=/bin/sh /start-transmission.sh